# recipes-bsp/swap/files/createswap.service
[Unit]
Description=Expand your memory

[Service]
Type=oneshot
ExecStart=/bin/sh -c "/usr/bin/createswap.sh"

[Install]
WantedBy=multi-user.target
